# daryrei-bot.service
[Unit]
Description=DaryRei Bot
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/daryrei-bot
ExecStart=/var/www/daryrei-bot/venv/bin/python /var/www/daryrei-bot/bot.py
Restart=always
RestartSec=10
Environment=PATH=/var/www/daryrei-bot/venv/bin
Environment=PYTHONPATH=/var/www/daryrei-bot
Environment=ENVIRONMENT=production

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=daryrei-bot

[Install]
WantedBy=multi-user.target
